# docker/s6-rc.d/idle-check/run
#!/command/execlineb -P
s6-setuidgid mountain
/bin/bash -c "while true; do /home/mountain/bin/idle-check.sh >> /home/mountain/idle-check.log 2>&1; sleep 600; done"
