# docker/s6-rc.d/idle-check/run
#!/command/execlineb -P
s6-setuidgid oopsbox
/bin/bash -c "while true; do /oopsbox/bin/idle-check.sh >> /oopsbox/idle-check.log 2>&1; sleep 600; done"
