# docker/s6-rc.d/agents-init/up
s6-setuidgid oopsbox /oopsbox/bin/agents-init.sh
