# docker/s6-rc.d/system-term/up
s6-setuidgid oopsbox /oopsbox/bin/system-term.sh start
