# docker/s6-rc.d/system-term/up
s6-setuidgid mountain /home/mountain/bin/system-term.sh start
